# pysimCoder-Web/conf.d/blackbox.conf
[program:blackbox]
command=blackbox
autorestart=true
