# pysimCoder-Web/conf.d/fluxbox.conf
[program:fluxbox]
command=blackbox
autorestart=true
