# pysimCoder-Web/conf.d/xfce4-terminal.conf
[program:xfce4-terminal]
command=xfce4-terminal
